# misc/systemd/user/up.service.d/override.conf
[Service]
Environment="NOTIFY_BIN=/home/henrywallace/go/bin/notify"
