# systemd/user/up.timer
[Unit]
Requires=up.service

[Timer]
Unit=up.service
OnUnitActiveSec=1h
AccuracySec=1m

[Install]
WantedBy=timers.target
